# systemd/wayan-uran.service
[Unit]
Description=Wayan Uran — backup / root-fix Claude Code agent (Telegram gateway)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=wayan
Group=wayan
WorkingDirectory=/opt/wayan-uran
EnvironmentFile=/etc/wayan-uran.env
Environment=PYTHONUNBUFFERED=1
# Claude Code is installed per-user; make sure its bin dir is on PATH.
Environment=PATH=/home/wayan/.local/bin:/usr/local/bin:/usr/bin:/bin
# Verify the home dirs the gateway/Claude need exist before starting.
ExecStartPre=/usr/bin/test -d /home/wayan/.config
ExecStartPre=/usr/bin/test -d /home/wayan/.claude
ExecStart=/opt/wayan-uran/venv/bin/python -m gateway
Restart=on-failure
RestartSec=5
# Allow the loop to finish the current poll/reply on stop or restart.
TimeoutStopSec=15
SyslogIdentifier=wayan-uran
# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=read-only
# Claude writes auth/state/cache under the home dir; the gateway writes the workspace.
ReadWritePaths=/home/wayan/.claude-lab /opt/wayan-uran /home/wayan/.claude /home/wayan/.config /home/wayan/.cache

[Install]
WantedBy=multi-user.target
